# all/99-usb-can-no-autosuspend.rules
# Udev rule to disable USB autosuspend for the CAN adapter (gs_usb)
# This prevents the adapter from being suspended during print idle moments
#
# Install: sudo cp 99-usb-can-no-autosuspend.rules /etc/udev/rules.d/
# Reload: sudo udevadm control --reload-rules
#
# The gs_usb driver (USB CAN adapter) can disconnect if USB autosuspend kicks in,
# causing "Couldn't shutdown device (err=-19)" errors and Klipper communication failures.

# Match the USB device directly when gs_usb driver binds
# idVendor=1d50, idProduct=606f is the Klipper CAN adapter (stm32f407xx)
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1d50", ATTR{idProduct}=="606f", ATTR{power/autosuspend}="-1", ATTR{power/control}="on"

# Fallback: Also set when can0 interface appears (handles edge cases)
ACTION=="add", SUBSYSTEM=="net", KERNEL=="can0", RUN+="/bin/sh -c 'echo -1 > /sys/bus/usb/devices/1-1/power/autosuspend 2>/dev/null; echo on > /sys/bus/usb/devices/1-1/power/control 2>/dev/null'"
